# bin/birthday_mail.cron
#! /bin/sh
cd $1
$2 -f ./symfony sfadvanced:send-birthday-mail
